# config_files/websockify.conf
[program:websockify]
command=/system/config_files/websockify.sh
autorestart=true
